# pinotams.service.example
[Unit]
Description=NOTAM Scraper
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/pinotams/bin/pinotams -s
Restart=on-abort

[Install]
WantedBy=multi-user.target
